# run/rc.local
#!/bin/sh -e

cd /root/datascience-snippets/run
./boot.sh >> /var/log/boot-script

exit 0
