# run/rc.local
#!/bin/sh

cd /root/datascience-snippets/run
./boot.sh >> /var/log/boot-script
